# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
loglevel=info

; ── FastAPI backend ─────────────────────────────────────────
[program:fastapi]
command=uvicorn api.main:app --host 0.0.0.0 --port 8000 --workers 1 --log-level warning
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
priority=1

; ── Streamlit UI ────────────────────────────────────────────
[program:streamlit]
command=streamlit run ui/app.py --server.port=7860 --server.address=0.0.0.0 --server.headless=true --browser.gatherUsageStats=false
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
priority=2